# server/bader-admin.service
[Unit]
Description=BADER Admin Panel (Streamlit)
After=network.target bader-api.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/bader-server
Environment="PATH=/opt/bader-server/venv/bin:/usr/bin"
ExecStart=/opt/bader-server/venv/bin/streamlit run admin_panel.py --server.port 8501 --server.address 0.0.0.0 --server.headless true
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
